# services/x-watcher/com.claude.x-check-playwright.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.claude.x-check-playwright</string>

    <key>ProgramArguments</key>
    <array>
        <string>{{HOME}}/.claude-mind/bin/x-check-playwright</string>
    </array>

    <key>StartInterval</key>
    <integer>900</integer>

    <key>StandardOutPath</key>
    <string>{{HOME}}/.claude-mind/system/logs/x-playwright.log</string>

    <key>StandardErrorPath</key>
    <string>{{HOME}}/.claude-mind/system/logs/x-playwright.log</string>

    <key>RunAtLoad</key>
    <false/>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>{{HOME}}</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:{{HOME}}/.local/bin</string>
    </dict>
</dict>
</plist>
