# services/meeting-check/com.claude.meeting-check.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.claude.meeting-check</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/claude/.claude-mind/bin/meeting-check</string>
    </array>
    <key>StartInterval</key>
    <integer>900</integer>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/claude/.claude-mind/logs/meeting-check.stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/claude/.claude-mind/logs/meeting-check.stderr.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/claude</string>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin:/Users/claude/.local/bin</string>
        <key>MIND_PATH</key>
        <string>/Users/claude/.claude-mind</string>
    </dict>
</dict>
</plist>
