# services/webhook-receiver/com.claude.webhook-receiver.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.claude.webhook-receiver</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/claude/Developer/samara-main/services/webhook-receiver/server.py</string>
        <string>--port</string>
        <string>8082</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/claude/.claude-mind/logs/webhook-receiver.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/claude/.claude-mind/logs/webhook-receiver.log</string>

    <key>WorkingDirectory</key>
    <string>/Users/claude/Developer/samara-main/services/webhook-receiver</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/claude</string>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:/Users/claude/.local/bin</string>
        <key>PYTHONPATH</key>
        <string>/Users/claude/Developer/samara-main/services/webhook-receiver</string>
    </dict>
</dict>
</plist>
