# services/wake-scheduler/com.claude.wake-adaptive.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.claude.wake-adaptive</string>

    <key>ProgramArguments</key>
    <array>
        <string>{{HOME}}/.claude-mind/bin/wake-adaptive</string>
    </array>

    <key>StartInterval</key>
    <integer>900</integer>

    <key>StandardOutPath</key>
    <string>{{HOME}}/.claude-mind/system/logs/wake-adaptive.log</string>

    <key>StandardErrorPath</key>
    <string>{{HOME}}/.claude-mind/system/logs/wake-adaptive.log</string>

    <key>RunAtLoad</key>
    <false/>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>{{HOME}}</string>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:{{HOME}}/.local/bin</string>
    </dict>
</dict>
</plist>
