# services/mcp-memory-bridge/com.claude.memory-bridge.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.claude.memory-bridge</string>
    <key>ProgramArguments</key>
    <array>
        <string>{{HOME}}/.local/bin/uv</string>
        <string>run</string>
        <string>--directory</string>
        <string>{{REPO_PATH}}/services/mcp-memory-bridge</string>
        <string>python</string>
        <string>server.py</string>
        <string>--http</string>
        <string>--port</string>
        <string>8765</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>{{HOME}}/.claude-mind/system/logs/memory-bridge.log</string>
    <key>StandardErrorPath</key>
    <string>{{HOME}}/.claude-mind/system/logs/memory-bridge.log</string>
    <key>WorkingDirectory</key>
    <string>{{REPO_PATH}}/services/mcp-memory-bridge</string>
</dict>
</plist>
